# etc/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
chmod = 0770
chown= npstack:nginx

[supervisord]
logfile = /var/log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = false
minfds = 1024
minprocs = 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[program:nginx]
command = /usr/sbin/nginx -c /etc/nginx/nginx.conf
stdout_events_enabled = true
stderr_events_enabled = true

[program:php-fpm]
command = /usr/sbin/php-fpm -c /etc/php.ini -y /etc/php-fpm.d/www.conf
stdout_events_enabled = true
stderr_events_enabled = true

[program:nginx-log]
command = tail -f /var/log/nginx.log
stdout_events_enabled = true
stderr_events_enabled = true

[program:php-fpm-log]
command = tail -f /var/log/php-fpm.log
stdout_events_enabled = true
stderr_events_enabled = true

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler